# .config/systemd/user/syncthing.service
[Unit]
Description=Syncthing

[Service]
Type=simple
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
